# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/target/linux/comcerto2000/base-files/etc/init.d/devices
#!/bin/sh /etc/rc.common

START=60

start() {
	major=`grep membuf /proc/devices | awk '{ print $1 }'`

	rm -f /dev/membuf

	if [ -n "$major" ]; then
		mknod /dev/membuf c $major 0
	fi
}
