# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/system/files/bootwrapper.init
#!/bin/sh /etc/rc.common

. /lib/functions.sh
. /lib/functions/network.sh

START=45
start() {
	/etc/init.d/firewall boot &
}

stop() {
	/etc/init.d/firewall stop
}
