# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/system/files/bootwrapper2.init
#!/bin/sh /etc/rc.common

. /lib/functions.sh
. /lib/functions/network.sh

START=89
STOP=10
start() {
	/etc/init.d/webfilter start &
}

stop() {
	/etc/init.d/webfilter stop
}
