# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/system/files/schedule.init
#!/bin/sh /etc/rc.common

#/usr/lib/webif/functions.sh
#/lib/config/uci.sh


reload (){
	uci commit schedule
}
