# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/system/files/systemstats.init
#!/bin/sh /etc/rc.common

. /lib/functions.sh
. /lib/functions/network.sh

START=10
start() {
uci_set_state systemstats  system "" system_state
}
