# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/utils/usbmode/files/usbmode.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=20

start() {
        /sbin/usbmode -s
}
