# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/c2k_gpio_reset/files/etc/init.d/gpio_reset
#!/bin/sh /etc/rc.common
# 2018-01-04: Ganesh Reddy <ganesh.reddy@nxp.com>
# Add c2k_gpio_reset module support.
# Copyright (C) 2008 OpenWrt.org

START=10

start() {
        insmod c2krv340_gpio_reset
}

stop(){
        rmmod c2krv340_gpio_reset
}

