# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/uwsgi/files/usr/bin/uwsgi-launcher
#!/bin/sh /etc/rc.common

export UWSGI_PLUGIN_DIR=/usr/lib/uwsgi/plugins

start() {
	uwsgi -m --ini /etc/uwsgi/jsonrpc.ini &
	uwsgi -m --ini /etc/uwsgi/blockpage.ini &
	uwsgi -m --ini /etc/uwsgi/upload.ini &
}

stop() {
	killall -9 uwsgi
}

