# src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/base-files/files/etc/init.d/umount
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

STOP=99
stop() {
	sync
	sleep 5
	/bin/umount -a -d -r
	ubidetach -d 0 > /dev/console
	sleep 1
}
